# deploy/supervisord.conf
[program:tawhiri-download]
directory=/srv
user=tawhiri
autostart=true
autorestart=true
command=nice ionice -c3 /srv/tawhiri2/bin/tawhiri-download daemon -w -f /var/log/tawhiri/download.log -e daniel@localhost -q
stopsignal=INT
stopwaitsecs=30

[program:tawhiri-v1]
directory=/srv
user=tawhiri
autostart=true
autorestart=true
command=/srv/tawhiri3/bin/gunicorn --config /srv/tawhiri3/gunicorn_cfg.py tawhiri.api:app
